# platform/roles/wildfly/templates/wildfly.service.j2
[Unit]
Description=WildFly application server
After=network.target

[Service]
Type=simple
User=wildfly
Group=wildfly
ExecStart=/opt/wildfly/bin/run.sh

[Install]
WantedBy=multi-user.target